# _03_hadoop/src/main/resources/supervisor.conf

[program:487-yarn-NODEMANAGER]
command=/opt/cloudera/cm-agent/bin/cm "launch" "--" "/opt/cloudera/cm-agent/service/yarn/yarn.sh" "nodemanager"
autostart=true
directory=/var/run/cloudera-scm-agent/process/487-yarn-NODEMANAGER
stdout_logfile=/var/run/cloudera-scm-agent/process/487-yarn-NODEMANAGER/logs/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/var/run/cloudera-scm-agent/process/487-yarn-NODEMANAGER/logs/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
environment=CDH_PARQUET_HOME='/usr/lib/parquet',YARN_NODEMANAGER_OPTS='-Xms1073741824 -Xmx1073741824 {{JAVA_GC_ARGS}} -Dlibrary.leveldbjni.path=/var/run/cloudera-scm-agent/process/487-yarn-NODEMANAGER -Dhadoop.event.appender=,EventCatcher -XX:+HeapDumpOnOutOfMemoryError -XX:HeapDumpPath=/tmp/yarn_yarn-NODEMANAGER-04bdfd07893754493f46020994940cdd_pid{{PID}}.hprof -XX:OnOutOfMemoryError=/opt/cloudera/cm-agent/service/common/killparent.sh',CDH_MR2_HOME='/usr/lib/hadoop-mapreduce',CMF_PACKAGE_DIR='/opt/cloudera/cm-agent/service',CM_STATUS_CODES='STATUS_NONE HDFS_DFS_DIR_NOT_EMPTY HBASE_TABLE_DISABLED HBASE_TABLE_ENABLED JOBTRACKER_IN_STANDBY_MODE YARN_RM_IN_STANDBY_MODE',MGMT_HOME='/opt/cloudera/cm',CGROUP_ROOT_CPU='/sys/fs/cgroup/cpu,cpuacct',ORACLE_HOME='/usr/share/oracle/instantclient',CDH_IMPALA_HOME='/usr/lib/impala',CDH_YARN_HOME='/usr/lib/hadoop-yarn',CGROUP_GROUP_CPU='',CDH_HDFS_HOME='/usr/lib/hadoop-hdfs',HOME='/var/lib/hadoop-yarn',HADOOP_CREDSTORE_PASSWORD='czyrh0zrxgpwgexl76jmt1qwv',PATH='/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin',HADOOP_CLIENT_CONF_DIR='/etc/hadoop/conf.cloudera.yarn',CDH_HUE_PLUGINS_HOME='/usr/lib/hadoop',PARCEL_DIRNAMES='CDH-7.1.3-1.cdh7.1.3.p0.4992530',CDH_HADOOP_BIN='/usr/bin/hadoop',KEYTRUSTEE_KP_HOME='/usr/share/keytrustee-keyprovider',CDH_KUDU_HOME='/usr/lib/kudu',CLOUDERA_ORACLE_CONNECTOR_JAR='/usr/share/java/oracle-connector-java.jar',CONF_DIR='/var/run/cloudera-scm-agent/process/487-yarn-NODEMANAGER',CDH_SQOOP2_HOME='/usr/lib/sqoop2',KEYTRUSTEE_SERVER_HOME='/usr/lib/keytrustee-server',JSVC_HOME='/usr/libexec/bigtop-utils',NM_LOCAL_DIRS='/mnt/data1/yarn/nm',HIVE_DEFAULT_XML='/etc/hive/conf.dist/hive-default.xml',CLOUDERA_POSTGRESQL_JDBC_JAR='/opt/cloudera/cm/lib/postgresql-42.1.4.jre7.jar',CDH_HBASE_FILESYSTEM_HOME='/usr/lib/hbase-filesystem',PARCELS_ROOT='/opt/cloudera/parcels',CDH_KMS_HOME='/usr/lib/hadoop-kms',CGROUP_GROUP_CPUACCT='',CDH_HBASE_HOME='/usr/lib/hbase',CDH_SQOOP_HOME='/usr/lib/sqoop',CDH_HCAT_HOME='/usr/lib/hcatalog',CDH_OOZIE_HOME='/usr/lib/oozie',IS_KERBERIZED='true',CDH_HIVE_HOME='/usr/lib/hive',CDH_ZOOKEEPER_HOME='/usr/lib/zookeeper',CDH_HUE_HOME='/usr/lib/hue',CGROUP_GROUP_DEVICES='system.slice/cloudera-scm-agent.service',CGROUP_ROOT_MEMORY='/sys/fs/cgroup/memory',CLOUDERA_MYSQL_CONNECTOR_JAR='/usr/share/java/mysql-connector-java.jar',CDH_SPARK_HOME='/usr/lib/spark',CGROUP_GROUP_MEMORY='',KRB5CCNAME='krb5cc_cldr',YARN_LOGFILE='hadoop-cmf-yarn-NODEMANAGER-lt-cdp-node1.log.out',CDH_MR1_HOME='/usr/lib/hadoop-0.20-mapreduce',TOMCAT_HOME='/usr/lib/bigtop-tomcat',CDH_PIG_HOME='/usr/lib/pig',CDH_SENTRY_HOME='/usr/lib/sentry',CDH_VERSION='7',CGROUP_GROUP_BLKIO='',CDH_TEZ_HOME='/usr/lib/tez',CGROUP_ROOT_BLKIO='/sys/fs/cgroup/blkio',CDH_CRUNCH_HOME='/usr/lib/crunch',CDH_LLAMA_HOME='/usr/lib/llama/',CDH_HTTPFS_HOME='/usr/lib/hadoop-httpfs',YARN_ROOT_LOGGER='INFO,RFA',CDH_HADOOP_HOME='/usr/lib/hadoop',CDH_HBASE_CONNECTORS_HOME='/usr/lib/hbase-connectors',CMF_SERVER_ROOT='/opt/cloudera/cm',LLAP_KEYTAB_LINK_PATH='/var/lib/hadoop-yarn/hive.keytab',WEBHCAT_DEFAULT_XML='/etc/hive-webhcat/conf.dist/webhcat-default.xml',CDH_KAFKA_HOME='/usr/lib/kafka',CGROUP_ROOT_DEVICES='/sys/fs/cgroup/devices',CDH_KNOX_HOME='/usr/lib/knox',YARN_LOG_DIR='/var/log/hadoop-yarn',CDH_HBASE_INDEXER_HOME='/usr/lib/hbase-solr',SCM_DEFINES_SCRIPTS='/opt/cloudera/parcels/CDH-7.1.3-1.cdh7.1.3.p0.4992530/meta/cdh_env.sh',CDH_SOLR_HOME='/usr/lib/solr',CGROUP_ROOT_CPUACCT='/sys/fs/cgroup/cpu,cpuacct',CDH_FLUME_HOME='/usr/lib/flume-ng'
user=yarn
group=hadoop
exitcodes=0,2
autorestart=unexpected
startsecs=20
startretries=3
stopwaitsecs=30
killasgroup=true
# rlimits
